# deploy/systemd/support-portal.service
[Unit]
Description=Customer Support Portal (Lab Demo)
After=network.target

[Service]
Type=simple
User=support-portal
WorkingDirectory=/opt/support-portal
Environment=SUPPORT_PORTAL_BASE=/opt/support-portal
ExecStart=/opt/support-portal/venv/bin/python -m support_portal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
