# demo_c2/systemd/demo-c2.service
[Unit]
Description=Demo C2 Listener (Lab Demo — port 4444)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/demo-c2/listener.py
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
